# intstream-X.X/lib/systemd/system/celeryworker.service
[Unit]
Description=celery worker service
After=network.target

[Service]
User=intstream
EnvironmentFile=/etc/intstream/environment.conf
ExecStart=/usr/share/intstream/.venv/bin/celery  --max-memory-per-child 24000 --workdir=/usr/share/intstream/backend -A backend.celery:celery_app worker -l info

[Install]
WantedBy=multi-user.target

